# udev/99-unplug-power.rules
ACTION=="change", SUBSYSTEM=="typec", RUN+="/dev/home/jpm/scripts/thinkpad/power-plug.sh"
