# systemd/user/unmute-audio.service
[Unit]
Description=Unmute audio for "Speaker" sinks

[Service]
Type=oneshot
ExecStart=/bin/echo '' >/dev/null
